# root/etc/services.d/code-server/run
#!/usr/bin/with-contenv bash

if [ -n "${PASSWORD}" ]; then
  AUTH="password"
else
  AUTH="none"
  echo "starting with no password"
fi

if [ -n "${SSL_KEY}" ] && [ -n "${SSL_CERT}" ]; then
	SSL_CONFIG=" --cert ${SSL_CERT} --cert-key ${SSL_KEY}"
else
	echo "Running insecurely on HTTP."
	SSL_CONFIG=""
fi

exec \
	s6-setuidgid abc \
		/usr/bin/code-server --port 8443 --user-data-dir /config/data --extensions-dir /config/extensions --disable-telemetry --disable-updates --auth "${AUTH}" "${SSL_CONFIG}" /config/workspace
