# root/etc/services.d/code-server/run
#!/usr/bin/with-contenv bash

if [ -n "${PASSWORD}" ]; then
  AUTH="password"
else
  AUTH="none"
  echo "starting with no password"
fi

if [ -n "${BASEPATH}" ]; then
  BASEPATH="--base-path ${BASEPATH}"
else
  BASEPATH=""
fi

exec \
	s6-setuidgid abc \
		/usr/bin/code-server \
			--port 8443 \
			--user-data-dir /config/data \
			--extensions-dir /config/extensions \
			--disable-telemetry \
			--disable-updates \
			--auth "${AUTH}" \
			${BASEPATH} \
			/config/workspace
