# root/etc/services.d/code-server/run
#!/usr/bin/with-contenv bash

cd /opt/code-server || exit

if [ -n "${PASSWORD}" ]; then
  PASSARG="-P ${PASSWORD}"
else
  PASSARG="-N"
  echo "starting with no password"
fi

exec \
	s6-setuidgid abc \
		/opt/code-server/code-server \
			-H \
			${PASSARG} \
			-d /config/data \
			-e /config/extensions \
			--disable-telemetry \
			/config/workspace
