# root/etc/s6-overlay/s6-rc.d/init-code-server/run
#!/usr/bin/with-contenv bash

mkdir -p "${HOME}"/{extensions,data,workspace,.ssh}

if [ -n "${SUDO_PASSWORD}" ] || [ -n "${SUDO_PASSWORD_HASH}" ]; then
    echo "setting up sudo access"
    if ! grep -q 'abc' /etc/sudoers; then
        echo "adding abc to sudoers"
        echo "abc ALL=(ALL:ALL) ALL" >> /etc/sudoers
    fi
    if [ -n "${SUDO_PASSWORD_HASH}" ]; then
        echo "setting sudo password using sudo password hash"
        sed -i "s|^abc:\!:|abc:${SUDO_PASSWORD_HASH}:|" /etc/shadow
    else
        echo "setting sudo password using SUDO_PASSWORD env var"
        echo -e "${SUDO_PASSWORD}\n${SUDO_PASSWORD}" | passwd abc
    fi
fi

[[ ! -f "${HOME}/.bashrc" ]] && \
    cp /root/.bashrc "${HOME}/.bashrc"
[[ ! -f "${HOME}/.profile" ]] && \
    cp /root/.profile "${HOME}/.profile"

# fix permissions (ignore contents of ${HOME}/workspace)
find "${HOME}" -path "${HOME}/workspace" -prune -o -exec chown abc:abc {} +
chown abc:abc "${HOME}/workspace"
chmod 700 "${HOME}/.ssh"
if [ -n "$(ls -A "${HOME}/.ssh")" ]; then
    chmod 600 "${HOME}"/.ssh/*
fi

if [ -n "${GIT_USER_NAME}" ]; then
    git config --global user.name "${GIT_USER_NAME}"
    git config --global user.email "${GIT_USER_EMAIL}"
fi

# echo $S3_ACCESS_KEY_ID:$S3_SECRET_ACCESS_KEY > /etc/passwd-s3fs
# chmod 600 /etc/passwd-s3fs
if [ -n "${S3_ACCESS_KEY_ID}" ]; then
    echo $S3_ACCESS_KEY_ID:$S3_SECRET_ACCESS_KEY > "${HOME}/.passwd-s3fs"
    chmod 600 "${HOME}/.passwd-s3fs"
    mkdir /mnt/s3
    mkdir -p /tmp/s3cache
    chmod 777 /tmp/s3cache
    s3_region_opt=()
    if [ -n "${S3_REGION}" ]; then
        s3_region_opt=(-o endpoint=${S3_REGION})
    fi
    s3fs ${S3_BUCKET} /mnt/s3 -o passwd_file="${HOME}/.passwd-s3fs" -o url=${S3_ENDPOINT} "${s3_region_opt[@]}" -o use_path_request_style -o uid=1000 -o gid=1000 -o umask=0002 -o allow_other -o mp_umask=0002 -o multireq_max=5 -o use_cache=/tmp/s3cache -o nonempty -o use_xattr -o enable_noobj_cache -o parallel_count=15 -o multipart_size=64 -o max_stat_cache_size=100000 -o enable_noobj_cache #-o dbglevel=info -o curldbg

    # add SSH file

    read_count=3
    interval_time=5
    ssh_file="/mnt/s3/.ssh/id_ed25519"
    for((i=1;i<=$read_count;i++))
    do
    if [ -f "$ssh_file" ]; then
        echo "ssh file exist ***"
        # eval "$(ssh-agent -s)"
        cp "$ssh_file" "${HOME}/.ssh/"
        chown abc:abc "${HOME}/.ssh"
        chmod 600 "${HOME}/.ssh/id_ed25519"
        # ssh-keyscan -t rsa github.com >> ~/.ssh/known_hosts
        # ssh-add "${HOME}/.ssh/id_ed25519"
        # ssh -T git@github.com

        echo 'eval "$(ssh-agent -s)"' >> "${HOME}/.bashrc"
        echo "ssh-add ${HOME}/.ssh/id_ed25519" >> "${HOME}/.bashrc"
        echo 'export SSH_PRIVATE_KEY="$(cat /mnt/s3/.ssh/id_ed25519)"' >> "${HOME}/.bashrc"

        # cat "${HOME}/.ssh/id_ed25519" | tr -d '\r' | ssh-add - > /dev/null
        # ssh-keyscan github.com >> "${HOME}/.ssh/known_hosts"
        # ssh-add "${HOME}/.ssh/id_ed25519"
        # ssh-keyscan -t rsa github.com >> "${HOME}/.ssh/known_hosts"
        # ssh -T git@github.com
        if [ -f "/mnt/s3/init-code-server.sh" ]; then
            echo "init-code-server.sh exist ***"
            chmod +x /mnt/s3/init-code-server.sh
            /mnt/s3/init-code-server.sh
        fi
        break
    else
        echo "ssh file not exist, wait $interval_time seconds"
        sleep $interval_time
    fi
    done
fi

# echo "installing extensions***"
# /app/code-server/bin/code-server  --extensions-dir "${HOME}/extensions"  --install-extension IronGeek.vscode-env
mv  -v /temp/extensions/* "${HOME}/extensions/"

CODE_BIN=/app/code-server/bin/code-server
EXT_DIR="${HOME}/extensions"
echo "[init-extensions] Starting runtime extension installation..." >&2
install_one() {
  i=0; until [ $i -ge 3 ]; do
    if "$CODE_BIN" --extensions-dir "$EXT_DIR" --install-extension "$1"; then
      echo "[init-extensions] Installed $1" >&2
      return 0
    fi
    i=$((i+1)); sleep $((2**i))
  done
  echo "[init-extensions] Failed to install $1" >&2
}
echo "[init-extensions] Extensions to install: $EXTENSIONS_RUNTIME" >&2
for ext in $EXTENSIONS_RUNTIME; do
  [ -n "$ext" ] && install_one "$ext"
done
echo "[init-extensions] Finished runtime extension installation." >&2


if [ -n "${TOKEN_NGROK}" ]; then
    echo "starting ngrok"
    /app/ngrok config add-authtoken ${TOKEN_NGROK}
    mkdir -p /app/logs
    /app/ngrok http --log=stdout 3000 > /app/logs/ngrok.log & 
    sleep 7 
    grep "started tunnel" /app/logs/ngrok.log | awk '{print $NF}' > "${HOME}/workspace/https_url.txt"
fi

# expose PORTS with cloudflare
IFS=',' read -ra PORTS <<< "${HTTPS_PORTS:-8443}"
for port in "${PORTS[@]}"; do
    cloudflared tunnel --url http://localhost:$port --logfile /app/logs/cloudflared_$port.log &
done

# Wait for the tunnels to be established and extract the URLs
sleep 7
for port in "${PORTS[@]}"; do
    cf_url=$(grep -oE 'https://[a-zA-Z0-9.-]+\.trycloudflare\.com' /app/logs/cloudflared_$port.log | head -n 1)
    if [ -n "$cf_url" ]; then
        echo "# PORT $port" >> "${HOME}/workspace/https_url.txt"
        echo "$cf_url" >> "${HOME}/workspace/https_url.txt"
    fi
done

welcome_message='echo "Welcome to code-server terminal! power by Erik Vargas"'
if ! grep -Fxq "$welcome_message" "${HOME}/.bashrc"; then
    echo "$welcome_message" >> "${HOME}/.bashrc"
fi
